# ai-native-linux-os/systemd/ai-shell.service
[Unit]
Description=AI Shell Assistant Daemon
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/ai-native-linux
ExecStart=/opt/ai-native-linux/venv/bin/python /opt/ai-native-linux/src/ai_shell/ai_shell_daemon.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target 